# backup-scripts/backup-system.timer
[Unit]
Description=Back up system every morning

[Timer]
OnCalendar=*-*-* 2:00:00
Persistent=true

[Install]
WantedBy=timers.target